# config/renupharm.crontab
0 1 * * * docker exec -t `docker ps --format "{{.Names}}" | grep renupharm | grep app | head -1` bin/rake renumpharm:remove_expired_listings >> /home/ubuntu/logs/cron.log 2>&1
0,10,20,30,40,50 * * * * docker exec -t `docker ps --format "{{.Names}}" | grep renupharm | grep app | head -1` bin/rake sunspot:reindex >> /home/ubuntu/logs/cron.log 2>&1
30 8 * * 0 docker exec `docker ps --format "{{.Names}}" | grep nginx` /opt/cert_renewal/cert_renewal.sh --renew-all >> /home/ubuntu/logs/cron.letsencrypt.log 2>&1
