# config/renupharm.crontab
CONTAINER=`docker ps --format "{{.Names}}" | grep renupharm | grep app | head -1`
8 * * * * docker exec -t $CONTAINER bin/rake renumpharm:remove_expired_listings > /home/ec2-user/logs/cron.log 2>&1
