# config/renupharm.crontab
0 1 * * * docker exec -t `docker ps --format "{{.Names}}" | grep renupharm | grep app | head -1` bin/rake renumpharm:remove_expired_listings > /home/ec2-user/logs/cron.log 2>&1
0,10,20,30,40,50 * * * * docker exec -t `docker ps --format "{{.Names}}" | grep renupharm | grep app | head -1` bin/rake sunspot:reindex > /home/ec2-user/logs/cron.log 2>&1
30 8 * * 0 docker exec -it `docker ps --format "{{.Names}}" | grep renupharm | grep web` /opt/cert_renewal/cert_renewal.sh --renew-all >> /var/log/cron.letsencrypt.log 2>&1
